# etc/systemd/system/conky.service.d/override.conf
[Service]
StartLimitIntervalSec=60
StartLimitBurst=5